# packaging/systemd/blockstead-update.service
[Unit]
Description=Install a newer Blockstead
Documentation=https://github.com/LordMalachi/blockstead
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/lib/blockstead/blockstead-update
# If systemd has to terminate the helper, convert any still-active attempt into
# an explicit unverified failure after the updater process group has stopped.
ExecStopPost=/usr/lib/blockstead/blockstead-update --finalize-interrupted
User=root
Group=root
UMask=0022

# This unit is deliberately not sandboxed the way blockstead.service is. It has
# to replace /opt/blockstead, run systemctl, and keep running while the
# dashboard it restarts goes down and comes back. Its safety comes from the
# helper independently requiring the requested hash to equal the current
# passing-main manifest, and from the installer restoring the previous version
# by itself when the new one is not healthy.
#
# Building the dashboard and migrating the database can take a while on a small
# home machine, so the timeout is generous rather than tight.
TimeoutStartSec=3600
TimeoutStopSec=1800
